# librespot/rootfs/etc/s6-overlay/s6-rc.d/avahi/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start avahi service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

#until [ -e /var/run/dbus/system_bus_socket ]; do
#  sleep 1s
#done

bashio::log.info "Starting avahi-daemon..."

avahi-daemon --no-drop-root -f /etc/avahi/avahi-daemon.conf
