# MMDVMHost/etc/services.d/php81-fpm/run
#!/usr/bin/with-contenv bashio
php-fpm81
