# MMDVMHost/rootfs/etc/s6-overlay/s6-rc.d/php81-fpm/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting php-fpm81..."

php-fpm81 -F

