# MMDVMHost/rootfs/etc/s6-overlay/s6-rc.d/mmdvmhost/run
#!/usr/bin/with-contenv bashio

exec s6-setuidgid mmdvm /run.sh

