# librespot/rootfs/etc/s6-overlay/s6-rc.d/librespot/run
#!/usr/bin/with-contenv bashio

declare -a CONFIG_PATH
declare -a LIBRESPOT_OPTS

CONFIG_PATH=/data/options.json

LIBRESPOT_OPTS=$(jq --raw-output ".librespotopts" $CONFIG_PATH)

mkdir -p /share/snapfifo

until /usr/sbin/avahi-daemon -c; do
  sleep 1s
done
sleep 10s

bashio::log.info  "Start Librespot..."
exec librespot ${LIBRESPOT_OPTS}
