# xserver/rootfs/etc/s6-overlay/s6-rc.d/xserver/run
#!/usr/bin/with-contenv bashio

# ==============================================================================
# Home Assistant Add-on: xserver
# Runs the X Server
# ==============================================================================

bashio::log.info "Starting X Server..."

# Run X Server
if [[ "$(bashio::config 'xinit_command')" == "" ]] ; then
    bashio::log.info "X command: /usr/bin/X :0 $(bashio::config 'x_options')"
    /usr/bin/X :0 $(bashio::config 'x_options')
else
    echo "$(bashio::config 'xinit_command')" > /usr/bin/xinit-command.sh
    chmod +x /usr/bin/xinit-command.sh
    bashio::log.info "xinit command: /usr/bin/xinit /usr/bin/xinit-command.sh -- :0 $(bashio::config 'x_options')"
    /usr/bin/xinit /usr/bin/xinit-command.sh -- :0 $(bashio::config 'x_options')
fi

