# xserver/rootfs/etc/s6-overlay/s6-rc.d/installer/up
#!/command/execlineb -P

with-contenv

foreground { /etc/s6-overlay/scripts/installer }
